# stations/_archive/session-channel-py/dist/com.operonlab.session-channel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd plist for session-channel (per-user agent).

  Install:
    cp dist/com.operonlab.session-channel.plist ~/Library/LaunchAgents/
    sed -i '' "s|__HOME__|$HOME|g" ~/Library/LaunchAgents/com.operonlab.session-channel.plist
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.operonlab.session-channel.plist
    launchctl kickstart -k gui/$(id -u)/com.operonlab.session-channel
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.operonlab.session-channel</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.session-channel/.venv/bin/uvicorn</string>
        <string>main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>10101</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/.session-channel</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SESSION_CHANNEL_HOME</key>
        <string>__HOME__/.session-channel</string>
        <key>SESSION_CHANNEL_PORT</key>
        <string>10101</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/.session-channel/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.session-channel/logs/stderr.log</string>
</dict>
</plist>
