# stations/session-archiver/launchd/com.joneshong.session-archiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joneshong.session-archiver</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/joneshong/workshop/stations/session-archiver/scripts/run-archiver.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>5</integer>
        <key>Minute</key><integer>15</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/joneshong/.claude/data/session-archiver/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/joneshong/.claude/data/session-archiver/launchd-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/Users/joneshong/.local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
